# k3os/overlay/etc/init.d/k3s-restarter-service
#!/sbin/openrc-run

supervisor=supervise-daemon
name="k3s-restarter-service"
command="/sbin/k3s-restarter"
command_args=">/var/log/k3s-restarter.log 2>&1"

output_log=/var/log/k3s-restarter.log
error_log=/var/log/k3s-restarter.log

pidfile="/var/run/k3s-restarter.pid"
respawn_delay=5
respawn_max=0
