# package/harvester-os/files/etc/systemd/system/rke2-agent.service.d/override.conf
[Service]
ExecStartPre=/usr/sbin/harv-update-rke2-server-url agent
