# k3os/overlay/etc/profile.d/harvester_console_bind.sh
#!/bin/bash

# bind <F12>
bind -x '"\e[24~":"harvester-console"'

# dashboard mode
export HARVESTER_DASHBOARD=true
export TTY=$(tty)