# package/harvester-os/files/etc/systemd/system/cos-setup-initramfs.service.d/override.conf
[Service]
ExecStartPre=-/usr/sbin/migration.sh
